# cookbooks/java/files/etc/profile.d/java.sh

export JAVA_HOME="/usr/java"
export PATH="${JAVA_HOME}/bin:${PATH}"
